# configs/supervisor/project_name.conf
[program:*project_name*]
command=/root/project/bin/start_gunicorn.sh
user=root
process-name=%(program_name)s
numproc=1
redirect_stderr=true
autostart=1
autorestart=1